# overlay/etc/systemd/system/cryptofaxpa-printer.path
[Unit]
Description=CryptoFaxPA Thermal Printer
[Path]
PathExists=/dev/usb/lp0
[Install]
WantedBy=multi-user.target
