# overlay/etc/systemd/system/wificonf.service
[Unit]
Description = Wi-Fi Configuration daemon

[Service]
Type = simple
Restart = always
ExecStart = /home/pi/wificonf -listen 192.168.90.1:80
